# moo-notify.service
[Unit]
Description=moo-notify - Desktop Notification Service
After=network-online.target graphical-session.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/node %h/code/moo-notify/index.js
Restart=on-failure
RestartSec=30

# Pass the display environment so xcowsay can open a window
Environment=DISPLAY=:0
PassEnvironment=DISPLAY DBUS_SESSION_BUS_ADDRESS XAUTHORITY XDG_RUNTIME_DIR HUNGERHUB_EMPLOYEE_ID HUNGERHUB_TOKEN

[Install]
WantedBy=default.target
